# templates/etc/systemd/system/kubelet.service.j2
#jinja2: trim_blocks:False
[Unit]
Description=Kubernetes Kubelet
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=network-online.target
After=containerd.service
Requires=network-online.target
Requires=containerd.service
Wants=network-online.target

[Service]
ExecStart={{ k8s_worker_bin_dir }}/kubelet \
{%- for setting in k8s_worker_kubelet_settings|sort %}
  --{{ setting }}{% if k8s_worker_kubelet_settings[setting] != "" %}={{ k8s_worker_kubelet_settings[setting] }}{% endif %}{% if not loop.last %} \{% endif %}
{%- endfor %}

Restart=always
RestartSec=10
WatchdogSec=30s

[Install]
WantedBy=multi-user.target
