# templates/etc/systemd/system/kubelet.service.j2
#jinja2: trim_blocks:False
[Unit]
Description=Kubernetes Kubelet
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=containerd.service
Requires=containerd.service

[Service]
ExecStart={{k8s_bin_dir}}/kubelet \
{%- for setting in k8s_worker_kubelet_settings|sort %}
  --{{setting}}={{k8s_worker_kubelet_settings[setting]}} {% if not loop.last %}\{% endif %}
{%- endfor %}

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
